# cron/crontab
0 1 * * * "${APPROOT}/run-task" wp-db-backup
0 2 * * * "${APPROOT}/run-task" lerenew
